# observability/etc/init.d/vlproxy
#!/sbin/openrc-run

command="/usr/local/vicky/bin/victorialogs"
command_args="-httpListenAddr=:9429 \
    -storageNode=o11y-1.vpn.pluralkit.net:9428 \
    -storageNode=o11y-2.vpn.pluralkit.net:9428 \
    -storageNode=o11y-3.vpn.pluralkit.net:9428"
output_log="/var/log/vlproxy.log"
error_log="/var/log/vlproxy.err"
command_background="yes"
pidfile="/run/vlproxy.pid"
