# observability/etc/init.d/vmalert-metrics
#!/sbin/openrc-run

command="/usr/local/vicky/bin/vmalert"
command_args="-rule=/usr/share/vicky/vmalert-metrics.yml \
	-datasource.url=http://metrics.svc.pluralkit.net:80/select/0/prometheus/ \
	-remoteWrite.url=http://metrics.svc.pluralkit.net:80/insert/0/prometheus/ \
	-remoteRead.url=http://metrics.svc.pluralkit.net:80/select/0/prometheus/ \
	-notifier.url=http://o11y-1.vpn.pluralkit.net:9093 \
	-notifier.url=http://o11y-2.vpn.pluralkit.net:9093 \
	-notifier.url=http://o11y-3.vpn.pluralkit.net:9093 \
"
output_log="/var/log/vmalert-metrics.log"
error_log="/var/log/vmalert-metrics.err"
command_background="yes"
pidfile="/run/vmalert-metrics.pid"
