# observability/etc/init.d/vminsert
#!/sbin/openrc-run

command="/usr/local/vicky/bin/vminsert"
command_args="-storageNode=o11y-1.vpn.pluralkit.net \
    -storageNode=o11y-2.vpn.pluralkit.net \
	-storageNode=o11y-3.vpn.pluralkit.net"
output_log="/var/log/vminsert.log"
error_log="/var/log/vminsert.err"
command_background="yes"
pidfile="/run/vminsert.pid"
