# observability/etc/init.d/vmselect
#!/sbin/openrc-run

command="/usr/local/vicky/bin/vmselect"
command_args="-storageNode=o11y-1.vpn.pluralkit.net \
    -storageNode=o11y-2.vpn.pluralkit.net \
	-storageNode=o11y-3.vpn.pluralkit.net"
output_log="/var/log/vmselect.log"
error_log="/var/log/vmselect.err"
command_background="yes"
pidfile="/run/vmselect.pid"
